# deploy/tiny-film-battery.service
[Unit]
Description=Tiny Film battery monitor
After=local-fs.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__PROJECT_ROOT__
EnvironmentFile=-__PROJECT_ROOT__/.env
ExecStart=__PROJECT_ROOT__/scripts/run_battery.sh
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
